# Api/ApiService/Scripts/tson-api-v1.conf
kill timeout 300

description "TSON Api Service v1.1"
author "John Lyon-Smith"

start on runlevel [2345]
stop on runlevel [06]

setuid ubuntu

script
  /home/ubuntu/bin/tson-api-v1.sh
end script
